# ops/com.mastermold.autopilot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  launchd template for the Autopilot paper daemon (macOS local hosting).

  Install (adjust the repo path first if yours differs):
    cp ops/com.mastermold.autopilot.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.mastermold.autopilot.plist
  Remove:
    launchctl unload ~/Library/LaunchAgents/com.mastermold.autopilot.plist
    rm ~/Library/LaunchAgents/com.mastermold.autopilot.plist

  KeepAlive restarts the daemon if it dies. The bot only trades while its mode
  is "paper" (set from the Autopilot panel or POST /api/autopilot); otherwise
  the process idles. Laptop sleep pauses it — a Zo worker is the real 24/7 home.
  Logs land in /tmp/mastermold-autopilot.log.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.mastermold.autopilot</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>cd "$HOME/Documents/GitHub/mastermold" &amp;&amp; bin/autopilot-daemon</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/tmp/mastermold-autopilot.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/mastermold-autopilot.log</string>
</dict>
</plist>
